# podkop/files/etc/hotplug.d/iface/95-podkop
[ "$ACTION" = "ifup" ] || exit 0

. /lib/functions.sh

logger $ACTION $INTERFACE $DEVICE
